# lnmt.service
[Unit]
Description=Linux Network Management Toolbox Background Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/lnmt --service
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
